# service/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS updater (oneshot)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/cloudflare-ddns
ExecStart=/opt/cloudflare-ddns/venv/bin/python ddns.py
